# sway/.config/systemd/user/dunst.service
[Unit]
Description=Dunst notification daemon

[Service]
ExecStartPre=test -n "$SWAYSOCK"
ExecStart=dunst

[Install]
WantedBy=sway-session.target
