# sway/.config/systemd/user/nm-applet.service
[Unit]
Description=Run nm-applet (systray)

[Service]
ExecStartPre=test -n "$SWAYSOCK"
ExecStart=nm-applet

[Install]
WantedBy=sway-session.target
